# scripts/cloud/openclaw-stack.service
[Unit]
Description=OpenClaw + Nanobot Stack
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Environment=HOME=/home/ubuntu
Environment=PATH=/home/ubuntu/.npm-global/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/ubuntu/cloud/start-stack.sh
ExecStop=/home/ubuntu/cloud/stop-stack.sh
Restart=always
RestartSec=5
TimeoutStopSec=30
KillMode=process

[Install]
WantedBy=multi-user.target